# psu-adp.target
[Unit]
Description=using power adapter
DefaultDependencies=no
StopWhenUnneeded=no
