# psu-bat.target
[Unit]
Description=using battery
DefaultDependencies=no
StopWhenUnneeded=no
